# supervisor.conf
[program:goal-server]
command=/opt/app/goal
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/opt/app/storage/logs/goal.log
